# medicine/scripts/medicine-bot.service
[Unit]
Description=Medicine Tracker WhatsApp Bot
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/medicine/server
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
